# back-end/infra/systemd-user/memos-backend.service
[Unit]
Description=Memos Backend (DevStack)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/devstack/workspace
ExecStart=%h/devstack/workspace/scripts/run-memos-backend.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
